# raspi34-esp133-uart-only/root-usb/etc___systemd___system___usb-ensure-mounted.service
[Unit]
Description=Ensure USB is mounted at /mnt/usb (read-only)

[Service]
Type=oneshot
ExecStart=/home/amane/root-usb/usb-ensure-mounted.sh
